# extra/game_server.conf
description "uWSGI server"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

chdir /var/www/ctf_game_server/
exec sudo uwsgi --ini uwsgi.ini
