# extra/game_server/ctf_game_server/ctf_game_server.service
[Unit]
Description=uWSGI instance to serve the ctf game server
After=network.target

[Service]
User=vagrant
Group=www-data
WorkingDirectory=/var/www/ctf_game_server
ExecStart=/usr/local/bin/uwsgi --ini ctf_game_server.ini

[Install]
WantedBy=multi-user.target